# recipes-bsp/usbotg-configs/hid-gamepad-usbotg-gadget-config/97-ustotg.rules
ACTION=="add", SUBSYSTEM=="udc", RUN+="/sbin/usbotg_hid_gamepad_config.sh restart"

