# recipes-bsp/usbotg-configs/acm-usbotg-gadget-config/97-ustotg.rules
ACTION=="add", SUBSYSTEM=="udc", RUN+="/sbin/usbotg_acm_config.sh restart"

